# chattering_fix.service
[Unit]
Description=Keyboard Chattering Fix service

[Service]
# Change ExecStart to the absolute path of the file, executing chattering_fix.sh
ExecStart=<absolute/path/to/file/chattering_fix.sh>

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target