# chattering_fix.service
[Unit]
Description=Keyboard Chattering Fix service

[Service]
Restart=on-failure
ExecStart=cd <path-to-repo-root> && sudo python3 -m src

[Install]
WantedBy=multi-user.target