# chattering_fix.service
[Unit]
Description=Keyboard Chattering Fix service

[Service]
Restart=on-failure
ExecStart=<path-to-script>

[Install]
WantedBy=multi-user.target